# systemd/trading-bot-stock-local.service
[Unit]
Description=Trading Bot Stock Paper Session (Local User Service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/workspaces/trading-bot
ExecStart=/workspaces/trading-bot/tools/run_stock_session.sh
Restart=on-failure
RestartSec=30
KillSignal=SIGINT
TimeoutStopSec=60

# Keep logs in the journal; the bot also writes stock_bot.log locally.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
