# systemd/nas/trading-bot-watchdog.timer
[Unit]
Description=Run Trading Bot NAS Health Watchdog Every 15 Minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
Unit=trading-bot-watchdog.service

[Install]
WantedBy=timers.target
