# systemd/nas/trading-bot-watchdog.service
[Unit]
Description=Trading Bot NAS Health Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/nas/trading-bot
ExecStart=/home/nas/trading-bot/tools/nas_health_check.sh /home/nas/trading-bot
User=nas
Group=nas
