# systemd/nas/trading-bot-train-hourly.service
[Unit]
Description=Trading Bot Hourly Model Training + Discord Report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/nas/trading-bot
ExecStart=/bin/bash /home/nas/trading-bot/tools/hourly_train_and_report.sh /home/nas/trading-bot
User=nas
Group=nas
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target